# systemd/autobackup.service
[Unit]
Description=AutoBackup - Automatic backup service
After=network.target

[Service]
Type=simple
ExecStart=/home/samo/Data/Programs/Back-up/AutoBackUp/venv/bin/python -m autobackup --daemon
WorkingDirectory=/home/samo/Data/Programs/Back-up/AutoBackUp
Restart=on-failure
RestartSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autobackup

[Install]
WantedBy=default.target
